# openjdk-11-rocky-8/profile.d/bash.sh
PS1='\u@\h: \W $ '

if [ -d ~/bin ]; then
  export PATH=~/bin:$PATH
fi
