# amazoncorretto-8/profile.d/bash.sh
PS1='\u@C8: \W $ '

if [ -d ~/bin ]; then
  export PATH=~/bin:$PATH
fi
