# amazoncorretto-8/profile.d/java.sh
export JAVA_HOME=/usr/lib/jvm/java-8-amazon-corretto
java -version
